# add/apache.ini
[program:apache2]
#command=/path/to/httpd -c "ErrorLog /dev/stdout" -DFOREGROUND
command=/usr/sbin/apachectl  -D "FOREGROUND"  -k start
redirect_stderr=true
autostart=true
autorestart=true
