# add/mythbackend.ini
[program:mythbackend]
command=/usr/bin/mythbackend
user=mythtv
environment=HOME="/home/mythtv",USER="mythtv"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autostart=true
autorestart=true
